# com.matthias.overtime-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.matthias.overtime-tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>-m</string>
        <string>src.main</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/overtime-tracker.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/overtime-tracker.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
</dict>
</plist>
